# build/files/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:geth]
command=/usr/local/bin/geth --config ./config.toml --datadir ./data  --cache 18000 --rpc.allow-unprotected-txs --txlookuplimit 0 %(ENV_EXTRA_OPTS)s
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=INT


[program:nginx]
command=nginx -c /etc/nginx/nginx.conf
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
